# build/supervisor/jupyter.conf
[program:jupyter]
user=vagrant
directory=/vagrant
command=xvfb-run -a /home/vagrant/anaconda/bin/jupyter notebook --no-browser --config=/vagrant/build/jupyter-config.py
autostart=true
stderr_logfile=/home/vagrant/logs/jupyter.log
stdout_logfile=/home/vagrant/logs/jupyter.log
